# binance-balance-notifier.service
[Unit]
Description=foo

[Service]
User=%I
ExecStart=/usr/bin/node /usr/bin/bbnotifier/index.js

[Install]
WantedBy=multi-user.target